# pacmatic@.timer
[Unit]
Description=Email list of pending updates to %I timer

[Timer]
OnCalendar=daily UTC

[Install]
WantedBy=multi-user.target
